# django_app/ops/systemd/family-app-backup.timer
[Unit]
Description=Dagelijkse Family PostgreSQL-back-up

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=family-app-backup.service

[Install]
WantedBy=timers.target
